# networking_module/etc/cron.d/xroad-metrics-networking-cron
# Default cronjob for X-Road Metrics Networking Visualizer module
# Runs the Networking Visualizer data preparation script daily
# Script is run using the default networking module settings file (/etc/xroad-metrics/networking/settings.yaml)
# To run a specific profile, append the profile name to the end of the command so that the correct configuration
# is loaded


SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# m   h  dom mon dow  user              command
  12  4   *   *   *   xroad-metrics     /usr/share/xroad-metrics/networking/networking-cron-entrypoint.sh
